# docker/app/supervisor/worker/default-queue-listener.ini
[program:default-queue-listener]
numprocs=1
command=php /var/www/html/artisan queue:listen --queue=default
autostart=true
autorestart=true
user=www-data
group=www-data
