# docker/app/supervisor/worker/queue.ini
[program:queue-listener]
numprocs=1
command=php /var/www/html/artisan queue:listen
autostart=true
autorestart=true
user=www-data
group=www-data
