# docker/app/supervisor/worker/recursive-queue-listener.ini
[program:recursive-queue-listener]
numprocs=1
command=php /var/www/html/artisan queue:listen --queue=recursive-queue
autostart=true
autorestart=true
user=www-data
group=www-data
